# feeds/packages/utils/pcmciautils/files/etc/hotplug.d/.svn/entries
10

dir
28306
file:///raid/repositories/openwrt/branches/packages_8.09/utils/pcmciautils/files/etc/hotplug.d
file:///raid/repositories/openwrt



2008-11-03T21:38:32.595583Z
13113
florian














e0eaf0b8-c020-4168-9010-202975f5e11d

pcmcia_socket
dir

pcmcia
dir

